# source/gfjardim.usb.automount/etc/udev/rules.d/99_persistent_usbstick.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEMS=="usb", DRIVERS=="usb-storage" RUN+="/usr/local/sbin/usb_mount"