# source/gfjardim.usb.automount/etc/udev/rules.d/99_persistent_usbdisks.rules
ACTION=="add", KERNEL=="sd*", ENV{DEVTYPE}=="partition" RUN+="/usr/local/sbin/usb_mount"
ACTION=="remove", ENV{ID_TYPE}=="disk" RUN+="/usr/local/sbin/usb_umount"